# crontab
0 0 * * * uv run --project /CROUStillantDatasets /CROUStillantDatasets/__main__.py  2>&1
